# feeds/luci/themes/loong-generic/root/etc/uci-defaults/luci-theme-loong-generic
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Loong_Generic=/luci-static/loong-generic
	commit luci
EOF
exit 0
